# systemd/check-quota.timer
[Unit]
Description=Check Claude API quota every 30 minutes
Requires=check-quota.service

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
Unit=check-quota.service

[Install]
WantedBy=timers.target
