# systemd/check-quota.service
[Unit]
Description=Check Claude API quota usage
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/home/clap-admin/cooperation-platform/resource-sharing
ExecStart=/usr/bin/python3 check-quota.py
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
