# systemd/resource-share-daily-aggregate.service
[Unit]
Description=Resource-Share Daily Aggregation
After=network.target

[Service]
Type=oneshot
User=clap-admin
Group=clap-users
WorkingDirectory=/home/clap-admin/cooperation-platform/resource-sharing
ExecStart=/usr/bin/python3 aggregate_daily.py
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
